# gnome-remote-desktop/gnome-remote-desktop.initd
#!/sbin/openrc-run

supervisor="supervise-daemon"

name="GNOME Remote Desktop"
description="GNOME Remote Desktop Service"

command="/usr/libexec/gnome-remote-desktop-daemon"
#command_args="--system"
command_user="gnome-remote-desktop"

depend() {
    need dbus
    after display-manager
    want net
}

