# gnome-remote-desktop/gnome-remote-desktop.initd
#!/sbin/openrc-run

supervisor="supervise-daemon"

description="GNOME Remote Desktop"
command="/usr/libexec/gnome-remote-desktop"
command_args="--system"

